# telegram-reminder-bot.service
[Unit]
Description=Telegram Reminder Bot with Gemini AI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/telegram_reminder_bot_project
ExecStart=/root/telegram_reminder_bot_project/.venv/bin/python /root/telegram_reminder_bot_project/run_bot.py
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitInterval=60

# Memory Management
MemoryHigh=500M
MemoryMax=800M
MemorySwapMax=0

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 