# telegram-reminder-bot.service
[Unit]
Description=Telegram Reminder Bot with Gemini AI
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-reminder-gemini
ExecStart=/home/ubuntu/telegram-reminder-gemini/venv/bin/python /home/ubuntu/telegram-reminder-gemini/bot.py
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitInterval=60

# Memory Management
MemoryHigh=500M
MemoryMax=800M
MemorySwapMax=0

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 